# docker/finetune/swift/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord/supervisord.log
childlogdir=/var/log/
pidfile=/var/run/supervisord.pid


[program:jupyter]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/bin/bash -c "/etc/csghub/start_jupyter.sh"
priority=100

[program:swift]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/bin/bash -c "/etc/csghub/start.sh"
priority=1

[program:mem-monitor]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command=/bin/bash -c "/etc/csghub/mem_monitor.sh"
priority=110